# deploy/server-brain/server-brain.service
[Unit]
Description=JBrain2 server-brain neural wall display (LAN, no auth)
After=network-online.target
Wants=network-online.target

[Service]
# Adjust User and WorkingDirectory to your checkout.
User=jbrain
WorkingDirectory=/opt/JBrain2/deploy/server-brain
ExecStart=/usr/bin/python3 serve.py
Environment=BRAIN_PORT=8800
Environment=BRAIN_POWER_MAX_W=90
# Pin to the LAN interface instead of 0.0.0.0 if you want to be strict:
# Environment=BRAIN_HOST=192.168.1.10
Restart=on-failure
RestartSec=3
# Hardening: this only reads /proc and /sys and serves static bytes.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
